# geo/deploy/systemd/perbug-geo.service
[Unit]
Description=Perbug Geo API Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/perbug/backend
EnvironmentFile=/etc/perbug/geo.env
ExecStart=/usr/bin/node dist/server/index.js
Restart=always
RestartSec=2
StandardOutput=append:/var/log/perbug/geo-api.log
StandardError=append:/var/log/perbug/geo-api.error.log

[Install]
WantedBy=multi-user.target
